# 06_AdvancedDockerfile/complex_app/crontab
*/1 * * * * /opt/bullshit_generator.sh
*/1 * * * * /opt/feedme.sh
